# home/.config/systemd/user/chachi-blueman-applet.service
[Unit]
Name=chachi-blueman-applet.service
PartOf=chachi-sway-session.target

[Service]
ExecStart=/usr/bin/blueman-applet
