# home/.config/systemd/user/chachi-proton-vpn.service
[Unit]
PartOf=chachi-sway-session.target
Requires=chachi-sway-nm-applet.service

[Service]
ExecStart=/usr/bin/protonvpn-app
