# home/.config/systemd/user/chachi-polkit-gnome.service
[Unit]
PartOf=chachi-sway-session.target

[Service]
ExecStart=%h/.bin/polkit-agent
