# home/.config/systemd/user/chachi-slack.service
[Unit]
PartOf=chachi-sway-session.target

[Service]
ExecStart=/usr/bin/slack --startup
