# home/.config/systemd/user/chachi-polkit-gnome.service
[Unit]
PartOf=chachi-sway-session.target

[Service]
ExecStart=/usr/libexec/polkit-gnome-authentication-agent-1
