# home/.config/systemd/user/chachi-rclone-nas.service
[Unit]
PartOf=chachi-sway-session.target

[Service]
ExecStart=/usr/bin/rclone mount --vfs-cache-mode writes nas:/external-data/pablo %h/NAS/
