# home/.config/systemd/user/chachi-rclone-proton-drive.service
[Unit]
PartOf=chachi-sway-session.target

[Service]
ExecStart=/usr/bin/rclone mount proton-drive:/ ~/ProtonDrive/
