# home/.config/systemd/user/chachi-rclone-google-drive.service
[Unit]
PartOf=chachi-sway-session.target

[Service]
ExecStart=/usr/bin/rclone mount google-drive:/ ~/GoogleDrive/
