# home/.config/systemd/user/chachi-cpupower-gui.service
[Unit]
PartOf=chachi-sway-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/cpupower-gui --balanced
