# home/.config/systemd/user/chachi-pasystray.service
[Unit]
Name=chachi-pasystray.service

[Service]
ExecStart=/usr/bin/pasystray
