# home/.config/systemd/user/chachi-ulauncher.service
[Unit]
PartOf=chachi-sway-session.target

[Service]
ExecStart=/usr/bin/ulauncher --hide-window
