# home/.config/systemd/user/chachi-insync.service
[Unit]
PartOf=chachi-sway-session.target

[Service]
ExecStart=/usr/bin/insync start --no-daemon