# home/.config/systemd/user/chachi-nm-applet.service
[Unit]
Name=chachi-nm-applet.service
PartOf=chachi-sway-session.target

[Service]
ExecStart=/usr/bin/nm-applet
